# meowifier.service
[Unit]
Description=Meowifier Web Application
After=network.target

[Service]
Type=simple
User=meowifier
WorkingDirectory=/home/meowifier/meowifier-web
Environment="PATH=/home/meowifier/meowifier-web/venv/bin"

ExecStart=/home/meowifier/meowifier-web/venv/bin/python app.py

Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
